# src/twobit-gitpoller.init.in
#! /bin/sh
### BEGIN INIT INFO
# Provides:          @PACKAGE@
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/stop @PACKAGE@
### END INIT INFO

LOG=@LOGDIR@/@PACKAGE@.log
NAME=@PACKAGE@
PID=@PIDFILE@
SERVER=$(which twistd)
TAC=@PACKAGE@

case "$1" in
    start)
        logger "${NAME}: starting"
        echo -n "Starting ${NAME} ... "
        twistd --umask=0022 -l "${LOG}" --pidfile "${PID}" twobit_poller --config=@CONFFILE@
        if [ $? -eq 0 ]; then
            echo "success"
        else
            echo "failed"
        fi
        ;;
    stop)
        logger "${NAME}: stopping"
        echo -n "Stopping ${NAME} ... "
        if [ -f ${PID} ]; then
            kill $(cat ${PID}) && echo "success."
        else
            echo "nothing to stop."
        fi
        ;;
    status)
        echo -n "${NAME}: "
        if [ -f ${PID} ]; then
            if ps -e | awk '{print $1}' | grep -q $(cat ${PID}); then
                echo "running"
            else
                echo "not running"
            fi
        else
            echo "no pid file, you're on your own."
        fi
        ;;
    *)
        logger "${NAME}: Invalid usage"
        echo "Usage: $0 {start|stop}"
        exit 1
        ;;
esac

exit 0
